# systemd_units/dumpdata.timer
[Unit]
Description=Dump dummy data

[Timer]
OnBootSec=900s
Unit=dumpdata.service

[Install]
WantedBy=multi-user.target
