# systemd_units/stateRedRecovery.path
[Unit]
Description=State Red Firware Download App
DefaultDependencies=false
OnFailure=path-fail-notifier@%n.service

[Path]
PathExists=/tmp/stateRedEnabled
Unit=stateRedRecovery.service

[Install]
WantedBy=multi-user.target
