# devspec/systemd_units/amixer.service
[Unit]
Description=amixer
Before=rdkbrowser.service

[Service]
ExecStart=/usr/bin/amixer.sh
Type=idle

[Install]
WantedBy=multi-user.target
